# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[program:backend]
command=bun /app/backend/src/app.ts
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopasgroup=true
killasgroup=true
; The app budgets up to 15s for graceful shutdown (drain sockets, flush logs,
; close DB). supervisor's default 10s SIGKILL escalation cut that short.
stopwaitsecs=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",BACKEND_PORT="%(ENV_BACKEND_PORT)s",FRONTEND_PORT="%(ENV_FRONTEND_PORT)s"
priority=10

[program:nginx]
command=/app/docker/wait-and-start-nginx.sh
autostart=true
autorestart=true
exitcodes=0
startretries=3
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[eventlistener:crashnotify]
command=sh -c "while read line; do echo \"FATAL: process entered FATAL state - $line\" >&2; echo RESULT 2; echo OK; done"
events=PROCESS_STATE_FATAL
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
